# microservices/torun-jug/src/main/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:logstash]
command=/docker-entrypoint.sh -f /logstash.conf

[program:microservice]
command=java -jar /app.jar
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0